# test-env/supervisor/ai-secretary.conf
; ============================================================
;  AI Secretary — supervisor 設定
;  systemd の代わりに Docker コンテナ内でサービスを管理する
;  各サービスはラッパースクリプト経由で起動する
; ============================================================

; --- VOICEVOX Mock (コンテナ内蔵版) ---
[program:voicevox-mock]
command=/etc/supervisor/conf.d/start-voicevox-mock.sh
user=aiagent
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/supervisor/voicevox-mock.log
stderr_logfile=/var/log/supervisor/voicevox-mock-err.log
environment=PYTHONUNBUFFERED="1"
priority=100

; --- OpenClaw Gateway ---
; 未インストールの場合は sleep infinity で待機
[program:openclaw]
command=/etc/supervisor/conf.d/start-openclaw.sh
user=aiagent
autostart=true
autorestart=false
startsecs=5
startretries=1
stdout_logfile=/var/log/supervisor/openclaw.log
stderr_logfile=/var/log/supervisor/openclaw-err.log
priority=200

; --- voice-bridge (ヘッドレスモード) ---
; 未セットアップの場合は sleep infinity で待機
[program:voice-bridge]
command=/etc/supervisor/conf.d/start-voice-bridge.sh
user=aiagent
autostart=true
autorestart=false
startsecs=5
startretries=1
stdout_logfile=/var/log/supervisor/voice-bridge.log
stderr_logfile=/var/log/supervisor/voice-bridge-err.log
priority=300
